# docker/sftp/ubuntu/etc/rc.local
#!/bin/bash
pwsh -File "/opt/powershell/posh-sftp/mount.ps1"
exit 0